# pika-kde-settings/etc/xdg/autostart/pika-xdg-user-dirs-update-autostart.desktop
[Desktop Entry]
Name=Update XDG user-dirs
Exec=xdg-user-dirs-update
Type=Application
Terminal=false
NoDisplay=true
